# mezzanine/project_template/deploy/crontab
*/5 * * * * %(user)s %(manage)s poll_twitter
