# mezzanine/project_template/deploy/docker/crontab
*/5 * * * * root python /mezzanine-project/manage.py poll_twitter
